# systemd/wfb-node-refresh.service
[Unit]
Description=VISR boot orchestrator
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
ExecStart=/opt/visr/scripts/wfb-node-refresh.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
