# systemd/wfb-collect.service
[Unit]
Description=VISR WFB Collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/visr
Environment="PATH=/opt/visr/venv/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStartPre=/bin/sleep 5
ExecStart=/opt/visr/venv/bin/python /opt/visr/wfb_collect.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
